# synapse-lite-profile-switcher.service
[Unit]
Description=Synapse Lite Profile Switcher (kdotool)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=PATH=%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin

ExecStart=/usr/bin/python3 %h/.local/share/synapse-lite/synapse_lite_profile_switcher.py

Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
